# tests/srv6/trex_supervisor.conf
[program:trex]
command=/var/log/v2.87/t-rex-64 -i --stl
directory=/var/log/v2.87
process_name=trex
stdout_logfile=/tmp/trex.out.log
stderr_logfile=/tmp/trex.err.log
redirect_stderr=false
autostart=false
autorestart=false
startsecs=1
numprocs=1
stopasgroup=true
killasgroup=true
