# tests/srv6/trex_agent_supervisor.conf
[program:trex_agent]
command=python /root/trex_agent.py
directory=/root/
process_name=trex_agent
environment=PYTHONPATH="/var/log/v2.87/automation/trex_control_plane/interactive"
stdout_logfile=/tmp/trex_agent.out.log
stderr_logfile=/tmp/trex_agent.err.log
redirect_stderr=false
autostart=false
autorestart=false
startsecs=1
numprocs=1
