# ansible/roles/connection_db/templates/servercfgd.service.j2
[Unit]
Description=servercfgd

[Service]
ExecStart=/usr/bin/env python3 {{ servercfgd_path }}/servercfgd.py
