# ansible/roles/sonic-common/templates/etc/systemd/system/docker.service.d/http-proxy.conf.j2
[Service]
Environment={% for proxy in docker_proxy_list %}"{{ proxy }}" {% endfor %}
