# godiva-utils/gnmi-test/gnmitest/cmd/gnmitest_demo/docker/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log                ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                   ; do not daemonize
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/                           ; where child log files will live

[program:gnmitest]
command = /gnmitest_demo  -service_port 55555 -fake_agent_name "targetA" -fake_agent_port 55556 -fake_agent_textproto "targetA.textproto" -alsologtostderr
startsecs=10
startretries=10
autorestart=true
