# ansible/roles/vm_set/templates/mux-simulator.service.j2
[Unit]
Description=mux simulator
After=network.target

[Service]
ExecStart=/usr/bin/env {{python_command}} {{ abs_root_path }}/mux_simulator.py {{ mux_simulator_port }} {{ vm_set_name }} -v
